# nonlinear/fs-overlay/etc/systemd/system/playground.service
[Unit]
Description=Nonlinear-Labs Playground
After=syslog.target network.target

[Service]
ExecStart=/nonlinear/playground/playground

[Install]
WantedBy=multi-user.target
