# nonlinear/fs-overlay-usbrescue/etc/systemd/system/connect-to-nllserver.service
[Unit]
Description=Runs the script to connect to NLL server for maintenance purposes
After=network.target
Wants=network.target 

[Service]
ExecStart=/usr/bin/connect-to-nllserver.sh

[Install]
WantedBy=multi-user.target
